# ansible/roles/usb-content/files/etc_udev_rules.d_automount.rules
ACTION=="add",KERNEL=="sda1",RUN+="/bin/mount -o sync,noexec,nodev,noatime,nodiratime,utf8 /dev/%k /media/usb0"
ACTION=="remove", KERNEL=="sda1", RUN+="/bin/umount /dev/%k"
